# systemd/revflow-disk-monitor.service
[Unit]
Description=RevFlow Disk Space Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/revflow-os/scripts/disk-space-monitor.sh
StandardOutput=journal
StandardError=journal
