# systemd-services/revdispatch-real-frontend.service
[Unit]
Description=RevDispatch Frontend (nginx static)
After=network.target nginx.service
Requires=nginx.service

StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

Restart=on-failure
RestartSec=10s
CPUAccounting=true
MemoryAccounting=true
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target

