# systemd-services/revflow-heatmap.service
[Unit]
Description=RevFlow Heatmap API
After=network.target postgresql.service

StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/revflow-heatmaps/backend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/shared-api-engine/.env
ExecStart=/usr/bin/python3 /opt/revflow-heatmaps/backend/main.py

Restart=on-failure
RestartSec=10s
CPUAccounting=true
MemoryAccounting=true
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target

