# systemd-services/revflow-power-prompts.service
[Unit]
Description=RevPrompt Unified - Power Prompts API
After=network.target postgresql.service

StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Environment="PYTHONPATH=/opt/shared-api-engine"
Type=simple
User=root
WorkingDirectory=/opt/revprompt-unified
EnvironmentFile=/opt/shared-api-engine/.env
ExecStart=/opt/revflow-os/shared-venv/bin/python -m uvicorn api:app --host 0.0.0.0 --port 8700
StandardOutput=append:/var/log/revflow/revprompt-unified.log
StandardError=append:/var/log/revflow/revprompt-unified-error.log

Restart=on-failure
RestartSec=10s
CPUAccounting=true
MemoryAccounting=true
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target
