# systemd-services/revflow-self-heal.service
[Unit]
Description=RevFlow Self-Healing Service
After=network.target docker.service nginx.service
Wants=docker.service

StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/root/revflow-self-heal.sh
StandardOutput=append:/var/log/revflow-self-heal.log
StandardError=append:/var/log/revflow-self-heal.log

Restart=on-failure
RestartSec=10s
CPUAccounting=true
MemoryAccounting=true
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target

