# systemd-services/revflow-revintel.service
[Unit]
Description=RevFlow RevIntel Webhook Service
After=network.target postgresql.service

StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/revflow-lead-engine/backend/revintel
Environment="PYTHONPATH=/opt/revflow-lead-engine"
EnvironmentFile=/opt/shared-api-engine/.env
ExecStart=/usr/bin/python3 main.py
StandardOutput=append:/opt/revflow-lead-engine/logs/revintel.log
StandardError=append:/opt/revflow-lead-engine/logs/revintel-error.log

Restart=on-failure
RestartSec=10s
CPUAccounting=true
MemoryAccounting=true
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target

