# systemd-services/revflow-citations-nap.service
[Unit]
Description=RevFlow Citations NAP API (with venv)
After=network.target

StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/revflow-citations
ExecStart=/usr/bin/python3 nap_api.py

Restart=on-failure
RestartSec=10s
CPUAccounting=true
MemoryAccounting=true
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target

