# systemd-services/revhome-uvicorn.service
[Unit]
Description=RevHome RevScore IQ - Real Assessment Engine
After=network.target

StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/revhome
Environment="PYTHONPATH=/opt/revhome"
EnvironmentFile=/opt/shared-api-engine/.env
ExecStart=/opt/revflow-os/shared-venv/bin/python -m uvicorn revhome_chat_api:app --host 0.0.0.0 --port 8108 --log-level info

Restart=on-failure
RestartSec=10s
CPUAccounting=true
MemoryAccounting=true
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target

