# systemd-services/revpublish-backend.service
[Unit]
Description=RevPublish Backend API v2.0 (Port 8550)
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/revflow-os/modules/revpublish/backend
EnvironmentFile=/etc/systemd/system/revpublish-backend.env
ExecStart=/opt/revflow-os/shared-venv/bin/uvicorn main:app --host 0.0.0.0 --port 8550
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
[Service]
Type=simple
User=root
WorkingDirectory=/opt/revflow-os/modules/revpublish/backend
Environment="POSTGRES_PASSWORD=revflow2026"
Environment="POSTGRES_USER=postgres"
Environment="POSTGRES_DB=revflow"
Environment="POSTGRES_HOST=localhost"
Environment="POSTGRES_PORT=5432"
Restart=always
RestartSec=10

# Resource Limits
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target
